# services/UsbGadget/Storage.service
[Translation]
File=QtopiaServices
Context=UsbGadget/Storage
[Service]
Actions="startUsbService(QString)"
Name[]="USB Storage Service"
Gadget=QUsbStorageGadget
