# services/Calendar.service
[Service]
Actions = raiseToday();newEvent();newEvent(QDataTime,QDataTime,QString,QString);addEvent(PimEvent);removeEvent(PimEvent);updateEvent(PimEvent);showEvent(Quuid);showEvent(QUuid,QDate)
Icon = DateBook
Name = Calendar
Name[ja] = ???
[newEvent()]
Icon = DateBook
Name = New event
Name[ja] = ???
[raiseToday()]
Icon = today
Name = Today's calendar
Name[ja] = ???
