# services/PhotoEdit.service
[Translation]
File=QtopiaServices
Context=Photo Edit
[Service]
Icon=PhotoEdit
Name[]=Photo Edit
