# services/SystemMessages.service
[Translation]
File=QtopiaServices
Context=Email
[Service]
Actions = "postMessage(QString,QString)"
Icon = email
Name[]=SystemMessages
