# services/TimeMonitor.service
[Service]
Actions = timeChanged(QString)
Multiple = 1
