# services/TimeMonitor.service
[Service]
Actions = timeChange(QString)
Multiple = 1
