# services/WebAccess.service
[Service]
Actions = openURL(QString);retrieveURL(QString);getURL(QString,QString)
